# platform/riscv64/reptilesim.service
[Unit]
Description=ReptileSim Kiosk (Godot 4, RISC-V)
After=network-online.target graphical.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/reptilesim
Environment=GODOT_VULKAN_DEVICE=0
ExecStart=/opt/reptilesim/platform/riscv64/run-reptilesim.sh
Restart=on-failure

[Install]
WantedBy=graphical.target
